# files/etc/supervisor/conf.d/mongod.conf
[program:mongod]
command=/usr/bin/mongod -f /etc/mongod.conf --keyFile <MONGO_LOCAL>/mongodb-keyfile --replSet <REPLICA_SET>
autostart=true
autorestart=true
startretries=3

# Redirect output so we can see it using "docker logs"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
